# misc/almost-ddns.service
[Unit]
Description=Almost-DDNS
After=syslog.target network.target


[Service]
#ExecStart=almost-ddns

[Install]
WantedBy=multi-user.target

